# app/supervisord.conf
; SuPurrMente — un contenedor, tres procesos, tres usuarios.
; supervisord corre como root (PID gestionado por tini) y lanza cada programa con su
; usuario sin privilegios. Así un oauth2-proxy comprometido (usuario 'oauth') no puede
; leer ni el token de Whisker ni el secreto de Gmail (600, propiedad de 'tracker').
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:datasette]
command=/app/run-datasette.sh
user=datasette
directory=/app
autorestart=true
startretries=10
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:tracker]
command=/app/run-tracker.sh
user=tracker
directory=/app/src
autorestart=true
startretries=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:oauth2-proxy]
command=/app/run-oauth2.sh
user=oauth
directory=/app
autorestart=true
startretries=10
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
